# nexus
#! /bin/sh
#
# $FreeBSD: head/devel/git/files/git_daemon.in 307344 2012-11-12 20:32:17Z wxs $
#

# PROVIDE: nexus
# REQUIRE: DAEMON
# BEFORE:  apache24
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable nexus:
#
#nexus_enable="YES"

. /etc/rc.subr

name="nexus"
rcvar="nexus_enable"

load_rc_config $name

# No nexus_user, nexus_group here; daemon(8) requires root
# to write the pidfile.

: ${nexus_enable:=NO}
: ${nexus_directory:=/usr/local/nexus/nexus}
: ${nexus_pidfile:=/var/run/nexus.pid}
: ${nexus_flags:=-p ${nexus_pidfile} -u nexus}

nexus_chdir=${nexus_directory}

# https://stackoverflow.com/questions/1527049/
join_by () { local IFS=":"; shift; echo "$*"; }
classpath=$(join_by : $(echo ${nexus_directory}/lib/boot/*.jar))

mainclass=org.sonatype.nexus.karaf.NexusMain
java=/usr/local/bin/java

command="/usr/sbin/daemon"
command_args="${java} -cp ${classpath} -XX:+UnlockDiagnosticVMOptions -XX:+UnsyncloadClass -XX:+LogVMOutput -XX:LogFile=../sonatype-work/nexus3/log/jvm.log -Dkaraf.base=. -Dkaraf.data=../sonatype-work/nexus3 -Dkaraf.etc=etc/karaf -Dkaraf.home=. -Djava.util.logging.config.file=etc/karaf/java.util.logging.properties -Djava.io.tmpdir=../sonatype-work/nexus3/tmp ${mainclass} >/dev/null 2>&1"

run_rc_command "$1"
